# .crontab
0 * * * * ~/work/dotfiles/updateconfig.sh
